# templates/generic.service.erb
[Unit]
Description=Prometheus <%= @name %>
Wants=basic.target <%= @service_dependencies %>
After=basic.target network.target <%= @service_dependencies %>

[Service]
User=<%= @user %>
Group=<%= @group %>
EnvironmentFile=/etc/sysconfig/<%= @name%>
ExecStart=<%= @exporter_bin_path %> $OPTIONS

ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target