# templates/node_exporter.service.erb
[Unit]
Description=Prometheus <%= @exporter_name %>
Wants=basic.target
After=basic.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
EnvironmentFile=/etc/sysconfig/<%= @exporter_name%>
ExecStart=/opt/<%= @exporter_name %>-<%= @version %>.<%= @os %>-<%= @arch %>/<%= @exporter_name %> $OPTION

ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target